# .config/systemd/user/desktop.target
[Unit]
Description=Desktop environment
Requires=dbus.socket

[Install]
