# .config/systemd/user/gnirehtet.service
[Unit]
Description=Start gnirehtet reverse tethering

[Service]
ExecStart=/usr/bin/gnirehtet autorun

[Install]
WantedBy=default.target
