# .config/systemd/user/megasync.service
[Unit]
Description=Megasync

[Service]
ExecStart=/usr/bin/megasync
Environment=DISPLAY=:0
Nice=15

[Install]
WantedBy=default.target
