# .config/systemd/user/update-flatpaks.service
[Unit]
Description=Update flatpaks

[Service]
Type=oneshot
ExecStart=%h/.scripts/update-flatpaks

[Install]
WantedBy=default.target
