# .config/systemd/user/update-flatpaks.timer
[Unit]
Description=Update flatpaks

[Timer]
OnCalendar=*-*-* 10:10:00
Persistent=true

[Install]
WantedBy=timers.target
