# oldroles/symlinks/files/config/systemd/user/update-dotfiles-submodules.timer
[Unit]
Description=Update submodules in ~/.dotfiles daily

[Timer]
OnCalendar=*-*-* 10:00:00
Persistent=true

[Install]
WantedBy=timers.target
